# base/docker-ui/root/etc/services.d/dockerui/run
#!/command/with-contenv bash

exec s6-setuidgid abc /env/bin/python /app/main.py
