# apps/docker-rollarr/root/crontab
#run python script every hour
0 * * * * $(which python) /rollarr/PrerollUpdate.py > /proc/1/fd/1 2>/proc/1/fd/2
#empty
