# mod/docker-mod-rclone/root/etc/s6-overlay/s6-rc.d/init-mod-rclone-install/up
/etc/s6-overlay/s6-rc.d/init-mod-rclone-install/run