# apps/homelabarr-jackett/root/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

COMMAND="/app/Jackett/jackett --NoUpdates"

exec \
  s6-setuidgid abc \
  ${COMMAND} \
  ${RUN_OPTS} >/dev/null 2>&1

