# apps/docker-local-persist/root/etc/s6-overlay/s6-rc.d/svc-local-persist/run
  #!/usr/bin/with-contenv bash

  echo "Starting DockServer Local-Persist Volume Plugin..."
  exec s6-setuidgid abc /usr/local/bin/local-persist
  EOF
