# apps/homelabarr-deluge/root/etc/services.d/deluge-web/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc \
  /usr/bin/deluge-web \
  -d -c /config >/dev/null 2>&1
